# files/certbot.service
[Unit]
Description=Letsencrypt renewal service

[Service]
Type=oneshot
ExecStart=/snap/bin/certbot renew --quiet --agree-tos
